# isp-monitor/isp-monitor.service
[Unit]
Description=job to gather speedtest-cli data

[Service]
User=pi
Type=oneshot
ExecStart=/home/pi/isp-monitor.client
