# isp-monitor/isp-monitor.timer
[Unit]
Description=timer to run the speed test data gathering

[Timer]
OnBootSec=5m
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target